# programs.ini
[program:symfony_command_scheduler_program]
command=bash -c 'sleep 60 && exec php /var/www/symfony/bin/console command'
autorestart=true